# ansible/roles/mariadb/templates/clean-keystone.cron
*/2 * * * * root /usr/bin/clean-keystone-tokens  >> /var/log/clean-keystone.log 2>&1
